# systemd/outproxy-eagle.service
# outproxy-eagle systemd unit.
#
# Installed automatically by the .deb package. For a manual install, place the
# binary at /usr/bin/outproxy-eagle, the config at /etc/outproxy-eagle/eagle.conf
# and the web assets at /usr/share/outproxy-eagle/web, then:
#   sudo cp systemd/outproxy-eagle.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now outproxy-eagle
#
# Manage:
#   systemctl status|stop|start|restart outproxy-eagle
#   journalctl -u outproxy-eagle -f
#
# State (eagle.log, state.json) lives under /var/lib/outproxy-eagle, provided
# by StateDirectory below.

[Unit]
Description=outproxy-eagle - native HTTP/SOCKS5 outproxy with web statistics
Documentation=https://github.com/freeacetone/outproxy-eagle
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/outproxy-eagle -c /etc/outproxy-eagle/eagle.conf
WorkingDirectory=/var/lib/outproxy-eagle
Restart=on-failure
RestartSec=2

# Run as an unprivileged dynamic user; listeners use ports >1024 so no extra
# capabilities are required.
DynamicUser=yes
StateDirectory=outproxy-eagle

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Persisted data (state.json / eagle.log) lives under /var/lib/outproxy-eagle.
ReadWritePaths=/var/lib/outproxy-eagle

[Install]
WantedBy=multi-user.target
