# packaging/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    if [ -d /run/systemd/system ]; then
        systemctl daemon-reload >/dev/null 2>&1 || true
        systemctl enable outproxy-eagle.service >/dev/null 2>&1 || true

        if [ -z "$2" ]; then
            # Fresh install ($2 empty): start with the shipped default config.
            systemctl start outproxy-eagle.service >/dev/null 2>&1 || true
            echo "outproxy-eagle installed and started."
            echo "  edit /etc/outproxy-eagle/eagle.conf, then: systemctl restart outproxy-eagle"
        else
            # Upgrade ($2 = old version): pick up the new binary/unit only if the
            # service was already running (don't resurrect one stopped on purpose).
            systemctl try-restart outproxy-eagle.service >/dev/null 2>&1 || true
            echo "outproxy-eagle upgraded."
        fi
    fi
fi

exit 0
